# .cutover-backup/com.bmo.restart-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bmo.restart-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bmo/CC4Me-BMO/scripts/restart-watcher.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/bmo/CC4Me-BMO/logs/restart-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bmo/CC4Me-BMO/logs/restart-watcher-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/bmo/CC4Me-BMO</string>
</dict>
</plist>
