# .cutover-backup/com.bmo.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.assistant.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/bmo/CC4Me-BMO/daemon/dist/core/main.js</string>
        <string>/Users/bmo/CC4Me-BMO</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bmo/CC4Me-BMO</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/bmo/CC4Me-BMO/logs/daemon-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bmo/CC4Me-BMO/logs/daemon-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bmo</string>
    </dict>
</dict>
</plist>
