# launchd/com.assistant.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.assistant.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <!-- CRITICAL: Must be bootstrap.js, NOT main.js -->
        <string>/Users/agent/KKit-R2/daemon/dist/bootstrap.js</string>
        <string>/Users/agent/KKit-R2</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/agent/KKit-R2</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/agent</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/agent/KKit-R2/logs/daemon-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/agent/KKit-R2/logs/daemon-stderr.log</string>
</dict>
</plist>
